# files/rootfs/etc/profile.d/senzing.sh
#!/bin/bash

export SENZING_ROOT=/opt/senzing
export PYTHONPATH=${SENZING_ROOT}/g2/python
export LD_LIBRARY_PATH=${SENZING_ROOT}/g2/lib
